# ._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._._openwrt_files_slipnet.init
#!/bin/sh /etc/rc.common
START=99
start() { echo "SlipNet started" }
stop()  { echo "SlipNet stopped" }
